# root/etc/udev/rules.d/60-bluetooth-event.rules
ACTION=="add", SUBSYSTEM=="bluetooth", ENV{HOME}="/home/oz", RUN+="/home/oz/.local/bin/bluetooth_refresh"
ACTION=="remove", SUBSYSTEM=="bluetooth", ENV{HOME}="/home/oz", RUN+="/home/oz/.local/bin/bluetooth_refresh"
